# environment/kr260/jailhouse/install/etc/profile.d/jailhouse_path.sh
export JAILHOUSE_DIR="/root/jailhouse"
export PATH="${JAILHOUSE_DIR}/tools/:$PATH"
export NON_ROOTCELL_DIR="${JAILHOUSE_DIR}/inmates/demos"
export NON_ROOTCELL_LINUX_DIR="${JAILHOUSE_DIR}/inmates/linux"