# environment/kria/jailhouse/install/etc/profile.d/jailhouse_path.sh
export PATH="/root/jailhouse/tools/:$PATH"
export JAILHOUSE_DIR="/root/jailhouse"
export NON_ROOTCELL_DIR="/root/non_rootcell"
